# service/post-start/run
#!/bin/dash

# Makes /etc/service/nginx/run only run once
sv once post-start
exit 0
